# deploy/linux/mt5-quotes.service
[Unit]
Description=MT5 Quote HTTP/TCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mt5_project
Environment=MT5_BIND_HOST=0.0.0.0
Environment=MT5_TCP_PORT=9527
Environment=MT5_HTTP_PORT=9528
ExecStart=/usr/bin/node src/start-server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
